# infra/systemd/erislite-agent.service
    [Unit]
Description=ErisLITE Basalt Agent
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=erislite
WorkingDirectory=/opt/erislite
EnvironmentFile=/opt/erislite/.env
ExecStart=/usr/bin/python3 -m agent.agent_loop
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target